# tools/systemd/clideck-backup.timer
[Unit]
Description=Daily CliDeck data backup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
